# backup.service
[Unit]
Description=Timed backups

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/backup_script.sh
